# etc/jabberd-c2s.conf.in
description	"jabberd Client-2-Server component"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec sudo -u jabber @bindir@/@jabberd_c2s_bin@
